# site-modules/profile/templates/otel/otel-collector.service.epp
[Unit]
Description=OpenTelemetry Collector for Claude Code Monitoring
Documentation=https://opentelemetry.io/docs/collector/
After=docker.service network.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=true
WorkingDirectory=<%= $otel_dir %>
ExecStart=/usr/bin/docker-compose up -d
ExecStop=/usr/bin/docker-compose down
ExecReload=/usr/bin/docker-compose restart
TimeoutStartSec=300
TimeoutStopSec=60

# Restart policy
Restart=on-failure
RestartSec=10s

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=<%= $otel_dir %>

[Install]
WantedBy=multi-user.target
